# .config/systemd/user/thunar.service
[Unit]
Description=thunar daemon

[Service]
Type=simple
ExecStart=/sbin/thunar --daemon
