# .config/systemd/user/udiskie.service
[Unit]
Description=Automount external devices

[Service]
ExecStart=/sbin/udiskie
Type=exec

[Install]
WantedBy=multi-user.target
